# systemd/repo-builder.service
[Unit]
Description=Build Arch repo

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/repo.systemd.env
WorkingDirectory=/srv/repo_builder
ExecStart=/srv/repo_builder/arch_repo_builder ${REPO_ARCH}.yaml
Restart=on-failure