# systemd/repo-updater-full@.timer
[Unit]
Description=Partially update Arch repo weekly

[Timer]
OnCalendar=Sun

[Install]
WantedBy=timers.target