# systemd/repo-builder.service
[Unit]
Description=Build Arch repo
After=network-online.target
Wants=network-online.target

[Service]
User=!WRITE_USER_HERE!
WorkingDirectory=/srv/repo_builder
ExecStart=/usr/bin/python scripts/build_daemon.py --arch !WRITE_ARCH_HERE!
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target