# systemd/repo-sync.timer
[Unit]
Description=Sync Arch repo hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target