# systemd/repo-builder.service
[Unit]
Description=Build Arch repo

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/repo.systemd.env
WorkingDirectory=$REPO_BUILD_DIR
ExecStart=$REPO_BUILD_DIR/arch_repo_builder ${REPO_ARCH}.yaml
Restart=on-failure